# Demo_ball_3/simple-jenkins-influx-sqlserver-steamlit-main/jenkins/build/crontab
# START CRON JOB
PATH=/usr/local/bin
* * * * * python3 /app/main.py
* * * * * python3 /app/main2.py
# END CRON JOB