# simple-jenkins-influx-sqlserver-steamlit-main/jenkins/build/crontab
# START CRON JOB
PATH=/usr/local/bin
0 * * * * python3 /app/main.py
# END CRON JOB
